# deploy/openclaude.service
[Unit]
Description=openClaude Autonomous Service
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/openClaude
Environment=NODE_ENV=production
Environment=OPENCLAUDE_HOST=127.0.0.1
Environment=OPENCLAUDE_PORT=18883
ExecStart=/usr/bin/env node %h/openClaude/src/index.js
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
